# scripts/com.tars.split-routes.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tars.split-routes</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/usr/local/bin/tars-split-routes.sh</string>
    </array>

    <!-- Run at boot, on resolv.conf change, and every 30s to catch DNS rotation -->
    <key>RunAtLoad</key>
    <true/>
    <key>StartInterval</key>
    <integer>30</integer>
    <key>WatchPaths</key>
    <array>
        <string>/etc/resolv.conf</string>
    </array>

    <key>StandardOutPath</key>
    <string>/tmp/tars-split-routes.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/tars-split-routes.log</string>
</dict>
</plist>
